# squid/config/etc/supervicord.d/10-squid-exporter.conf
[program:squid_exporter]
command=/opt/squid-exporter -listen :8000
autorestart=true
